# containers/kafka/supervisor/kafka.conf
[program:kafka]
command=/usr/bin/start-kafka.sh
redirect_stderr=true
autostart=true
autorestart=true